# src/DEBIAN/postinst
#!/bin/bash
update-rc.d sleeponlan defaults
/etc/init.d/sleeponlan start
